# applications/pbase-node-solid-server/rpmbuild/SOURCES/pbase-node-solid-server-1.0/usr/local/pbase-data/pbase-node-solid-server/etc-systemd-system/solid.service
[Unit]
Description=solid - Social Linked Data
Documentation=https://solid.inrupt.com/docs/
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/yourhost.example.org
ExecStart=/usr/bin/solid start
Restart=on-failure

[Install]
WantedBy=multi-user.target
